# laravel-invest-worker.conf
[program:laravel-invest-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/laravel-invest/artisan queue:work redis --sleep=1 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=7
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/var/www/html/laravel-invest/storage/logs/worker.log
stopwaitsecs=3600

[program:laravel-invest-ml-scoring-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/laravel-invest/artisan queue:work redis --queue=ml-scoring --sleep=1 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=3
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/var/www/html/laravel-invest/storage/logs/ml-scoring-worker.log
stopwaitsecs=3600

[program:laravel-invest-ml-scoring-catchup-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/laravel-invest/artisan queue:work redis --queue=ml-scoring-catchup --sleep=1 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=1
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/var/www/html/laravel-invest/storage/logs/ml-scoring-catchup-worker.log
stopwaitsecs=3600

[program:laravel-invest-reverb]
process_name=%(program_name)s
command=php /var/www/html/laravel-invest/artisan reverb:start
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/laravel-invest/storage/logs/reverb.log
stopwaitsecs=10

[program:laravel-invest-backtest-b]
process_name=%(program_name)s
command=/var/www/html/laravel-invest/scripts/continuous-back/b-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=1
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-d]
process_name=%(program_name)s
command=/var/www/html/laravel-invest/scripts/continuous-back/d-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=1
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-e]
process_name=%(program_name)s
command=/var/www/html/laravel-invest/scripts/continuous-back/e-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=1
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-f]
process_name=%(program_name)s
command=/var/www/html/laravel-invest/scripts/continuous-back/f-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=1
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-h]
process_name=%(program_name)s
command=/var/www/html/laravel-invest/scripts/continuous-back/h-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=1
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-i]
process_name=%(program_name)s
command=/var/www/html/laravel-invest/scripts/continuous-back/i-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=1
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-j]
process_name=%(program_name)s
command=/var/www/html/laravel-invest/scripts/continuous-back/j-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=1
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-k]
process_name=%(program_name)s
command=/var/www/html/laravel-invest/scripts/continuous-back/k-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=1
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-l]
process_name=%(program_name)s
command=/var/www/html/laravel-invest/scripts/continuous-back/l-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=1
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-n]
process_name=%(program_name)s
command=/var/www/html/laravel-invest/scripts/continuous-back/n-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=1
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-o]
process_name=%(program_name)s
command=/var/www/html/laravel-invest/scripts/continuous-back/o-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=1
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-p]
process_name=%(program_name)s
command=/var/www/html/laravel-invest/scripts/continuous-back/p-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=1
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-p]
process_name=%(program_name)s
command=/var/www/html/laravel-invest/scripts/continuous-back/p-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=1
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-q]
process_name=%(program_name)s
command=/var/www/html/laravel-invest/scripts/continuous-back/q-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=1
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-bar-stream]
process_name=%(program_name)s
command=/var/www/html/laravel-invest/scripts/log-bar-stream.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=1
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=15

[program:laravel-invest-pipeline-watcher]
process_name=%(program_name)s
command=/var/www/html/laravel-invest/scripts/log-pipeline-watcher.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=1
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=10

[program:laravel-invest-ml-scoring-daemon]
process_name=%(program_name)s
command=php /var/www/html/laravel-invest/artisan ml:scoring-daemon
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=1
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/var/www/html/laravel-invest/storage/logs/ml-scoring-daemon.log
stopwaitsecs=30

[program:laravel-invest-backtest-a]
process_name=%(program_name)s
command=/var/www/html/laravel-invest/scripts/continuous-back/a-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=1
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-c]
process_name=%(program_name)s
command=/var/www/html/laravel-invest/scripts/continuous-back/c-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=1
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-g]
process_name=%(program_name)s
command=/var/www/html/laravel-invest/scripts/continuous-back/g-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=1
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-m]
process_name=%(program_name)s
command=/var/www/html/laravel-invest/scripts/continuous-back/m-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=1
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:trading-realtime-watch]
process_name=%(program_name)s
command=php /var/www/html/laravel-invest/artisan trading:realtime-watch
directory=/var/www/html/laravel-invest
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/var/www/html/laravel-invest/storage/logs/trading-realtime-watch.log
stopwaitsecs=10

[program:laravel-invest-backtest-r-workers]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/laravel-invest/artisan queue:work redis --queue=default --sleep=1 --tries=3 --max-time=7200
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pnovack
numprocs=12
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/bin:/usr/bin:/bin",HOME="/var/www/html"
redirect_stderr=true
stdout_logfile=/var/www/html/laravel-invest/storage/logs/backtest-r-workers.log
stopwaitsecs=3600