# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0770
chown=root:www-data

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:apache]
command=/usr/sbin/apache2ctl -D FOREGROUND
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:laravel-invest-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/laravel-invest/artisan queue:work redis --sleep=1 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=7
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",HOME="/root"
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stopwaitsecs=3600

[program:laravel-invest-ml-scoring-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/laravel-invest/artisan queue:work redis --queue=ml-scoring --sleep=1 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=3
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",HOME="/root"
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stopwaitsecs=3600

[program:laravel-invest-ml-scoring-catchup-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/laravel-invest/artisan queue:work redis --queue=ml-scoring-catchup --sleep=1 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=1
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",HOME="/root"
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stopwaitsecs=3600

[program:laravel-invest-reverb]
process_name=%(program_name)s
command=php /var/www/html/laravel-invest/artisan reverb:start
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",HOME="/root"
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:laravel-invest-scheduler]
process_name=%(program_name)s
command=php /var/www/html/laravel-invest/artisan schedule:work
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",HOME="/root"
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:laravel-invest-vite]
process_name=%(program_name)s
command=npm run dev -- --host 0.0.0.0
directory=/var/www/html/laravel-invest
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",HOME="/root",NODE_ENV="development"
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:laravel-invest-bar-stream]
process_name=%(program_name)s
command=bash /var/www/html/laravel-invest/scripts/log-bar-stream.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",HOME="/root"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=15

[program:laravel-invest-pipeline-watcher]
process_name=%(program_name)s
command=bash /var/www/html/laravel-invest/scripts/log-pipeline-watcher.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",HOME="/root"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=10

[program:laravel-invest-ml-scoring-daemon]
process_name=%(program_name)s
command=php /var/www/html/laravel-invest/artisan ml:scoring-daemon
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",HOME="/root"
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

; ── Continuous Backtests (pipelines b–s) ──
[program:laravel-invest-backtest-b]
process_name=%(program_name)s
command=bash /var/www/html/laravel-invest/scripts/continuous-back/b-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",HOME="/root"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-d]
process_name=%(program_name)s
command=bash /var/www/html/laravel-invest/scripts/continuous-back/d-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",HOME="/root"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-e]
process_name=%(program_name)s
command=bash /var/www/html/laravel-invest/scripts/continuous-back/e-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",HOME="/root"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-f]
process_name=%(program_name)s
command=bash /var/www/html/laravel-invest/scripts/continuous-back/f-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",HOME="/root"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-h]
process_name=%(program_name)s
command=bash /var/www/html/laravel-invest/scripts/continuous-back/h-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",HOME="/root"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-i]
process_name=%(program_name)s
command=bash /var/www/html/laravel-invest/scripts/continuous-back/i-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",HOME="/root"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-j]
process_name=%(program_name)s
command=bash /var/www/html/laravel-invest/scripts/continuous-back/j-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",HOME="/root"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-k]
process_name=%(program_name)s
command=bash /var/www/html/laravel-invest/scripts/continuous-back/k-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",HOME="/root"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-l]
process_name=%(program_name)s
command=bash /var/www/html/laravel-invest/scripts/continuous-back/l-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",HOME="/root"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-n]
process_name=%(program_name)s
command=bash /var/www/html/laravel-invest/scripts/continuous-back/n-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",HOME="/root"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-o]
process_name=%(program_name)s
command=bash /var/www/html/laravel-invest/scripts/continuous-back/o-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",HOME="/root"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-p]
process_name=%(program_name)s
command=bash /var/www/html/laravel-invest/scripts/continuous-back/p-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",HOME="/root"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30

[program:laravel-invest-backtest-q]
process_name=%(program_name)s
command=bash /var/www/html/laravel-invest/scripts/continuous-back/q-backtest_comparison.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
environment=PATH="/var/www/html/laravel-invest/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",HOME="/root"
redirect_stderr=true
stdout_logfile=/dev/null
stopwaitsecs=30
